# server/rootfs/etc/services.d/api/run
#!/command/with-contenv bash
set -euo pipefail

exec /server/scripts/run_api.sh

